# dcgm-exporter/dcgm-exporter.service
[Unit]
Description=Systemd DCGM-Exporter Service Unit File

[Service]
ExecStart=dcgm-exporter

[Install]
WantedBy=multi-user.target